# Shares/etc/profile.d/grails.sh
export GRAILS_HOME=/srv/grails
export GRAILS_OPTS="-server -noverify -Xmx1G -Xms1G -XX:MaxPermSize=512m -XX:PermSize=512m -XX:+AggressiveOpts -Dfile.encoding=UTF-8"
export PATH=$PATH:$GRAILS_HOME/bin